# Dot Files/.config/systemd/user/change_wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=/home/angad/Projects/Github/Linux-Presets/Dot Files/.config/hypr/hyprpaper.sh
